# deployment/storagehub-backend.service
[Unit]
Description=StorageHub backend (FastAPI / uvicorn)
After=network.target mysql.service mariadb.service
Wants=network-online.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_USER}}
WorkingDirectory={{BACKEND_DIR}}
# Loads {{BACKEND_DIR}}/.env via pydantic-settings (WorkingDirectory)
ExecStart={{BACKEND_DIR}}/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port {{BACKEND_PORT}} --workers 2
Restart=always
RestartSec=3
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
